# deployment/oracle/wakili-api.service
[Unit]
Description=Wakili FastAPI Backend
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/wakili-backend
EnvironmentFile=/opt/wakili-backend/.env
ExecStart=/opt/wakili-backend/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target